# pkg/debian/service.upstart
# rackspace-monitoring-poller

start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [016]

respawn

script
    . /etc/default/rackspace-monitoring-poller
    if [ "x$ENABLED" = xtrue ]; then
        exec /usr/bin/rackspace-monitoring-poller serve --config $CONFIG_FILE $POLLER_SERVE_OPTS
    else
        echo "Not enabled"
    fi
end script